# client/agent/cert-agent.service
[Unit]
Description=Certificate Agent – manages TLS cert lifecycle for this node
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/cert-agent/agent.env
ExecStart=/usr/bin/python3 -m agent
WorkingDirectory=/opt/cert-agent
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/cert-agent /etc/nginx/certs
ProtectHome=true

[Install]
WantedBy=multi-user.target
